# packaging/scripts/preremove.sh
systemctl stop prometheus-downsampler.service
systemctl disable prometheus-downsampler.service
